# webspaced.service
[Unit]
Description=webspace-ng daemon

[Service]
Type=simple
Exec=/usr/local/bin/webspaced

[Install]
WantedBy=multi-user.target
